# aaa/etc/systemd/system/pip_install_aaa.path
[Unit]
Description=Monitor /var/www/aaa/aaa/pip.requirements for changes.
After=systemd-networkd-wait-online.service
Wants=systemd-networkd-wait-online.service
Wants=pip_install.service

[Path]
PathChanged=/var/www/aaa/aaa/pip.requirements

[Install]
WantedBy=multi-user.target