# uifng/etc/systemd/system/yarn_install.path
[Unit]
Description=Monitor /var/www/ui-frontend-ng/package.json for changes.
After=systemd-networkd-wait-online.service
Wants=systemd-networkd-wait-online.service
Wants=yarn_install.service

[Path]
PathChanged=/var/www/ui-frontend-ng/package.json

[Install]
WantedBy=multi-user.target