# api-collaboration/etc/systemd/system/pip_install_api.path
[Unit]
Description=Monitor /var/www/api/api/pip.requirements for changes.

[Path]
PathModified=/var/www/api/api/pip.requirements
Unit=pip_install_api.service

[Install]
WantedBy=multi-user.target

