# uib/etc/systemd/system/pip_install_ui.path
[Unit]
Description=Monitor /var/www/ui-backend/backend/pip.requirements for changes.
After=systemd-networkd-wait-online.service
Wants=systemd-networkd-wait-online.service
Wants=pip_install_ui.service

[Path]
PathChanged=/var/www/ui-backend/backend/pip.requirements

[Install]
WantedBy=multi-user.target